# Examples/BasicService/fake.service.d/hellomars.conf
[Service]
ExecStart=/bin/echo "hello mars"
